# mcollective/etc/cron.d/hsec_puppet_versions
# run puppet versions and owners
23 7 * * * root /usr/bin/perl /localservices/sbin/puppet_versions.pl > /dev/null 2>&1
45 8 * * * root /usr/bin/perl /localservices/sbin/puppet_owners.pl   > /dev/null 2>&1
