# mcollective/etc/cron.d/hsec_puppet_versions
#
# run puppet versions
55 7 * * * MY_USER_CONFIGURED_IN_MCOLLECTIVE /usr/bin/perl /localservices/sbin/puppet_versions.pl > /dev/null 2>&1
